# data/instances/ksdemo/etc/supervisor.ini
[program:ksdemo-ks]
directory=/data/instances/ksdemo/etc
command=/data/software/knowledgestore/bin/ksd-service.sh /data/instances/ksdemo/etc/ks.ttl /data/instances/ksdemo/etc/
user=root
autostart=True
autorestart=unexpected
redirect_stderr=True
priority=2
environment=JAVA_HOME="/usr/lib/jvm/java-8-oracle/jre/"

[program:ksdemo-v]
directory=/data/instances/ksdemo/etc
command=/data/software/virtuoso/bin/virtuoso-t -f -c /data/instances/ksdemo/etc/virtuoso.ini
user=root
autostart=True
autorestart=unexpected
redirect_stderr=True
priority=2

[program:sshd]
command=/usr/sbin/sshd -D
autostart=True
user=root
autorestart=unexpected
redirect_stderr=True
priority=2